# config/crontab
# Run scraper based on configured interval
* * * * * /app/scripts/run_scraper.sh
